# backend/Supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=apache2-foreground
autostart=true
autorestart=true

[program:reverb]
command=php artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true